# target_files/usr/lib/systemd/system/run_customizations.service
[Unit]
Description=Run customizations service

[Service]
ExecStart=/usr/local/sbin/run_customizations.sh

[Install]
WantedBy=multi-user.target
